# deploy/systemd/spx-backtest-prod.service
[Unit]
Description=SPX Backtest Prod UI
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/spx-data-collector
EnvironmentFile=/opt/spx-data-collector/.env
ExecStart=/opt/spx-data-collector/.venv/bin/python -m spx_collector.backtest_prod --host 127.0.0.1 --port 8789
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
