# deploy/systemd/spx-sqlite-backup.service
[Unit]
Description=Backup SPX SQLite database

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/SPX-Data-Collector
EnvironmentFile=/home/ubuntu/SPX-Data-Collector/.env
ExecStart=/home/ubuntu/SPX-Data-Collector/scripts/backup_sqlite.sh
