# deploy/systemd/spx-sqlite-backup.service
[Unit]
Description=Backup SPX SQLite database

[Service]
Type=oneshot
User=appuser
Group=appuser
WorkingDirectory=/opt/spx-data-collector
EnvironmentFile=/opt/spx-data-collector/.env
ExecStart=/opt/spx-data-collector/scripts/backup_sqlite.sh
